# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /code/uwsgi.ini